# Library/LaunchAgents/com.__username__.homebrewupdate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.__username__.homebrewupdate.plist</string> <!-- UPDATE file name -->

  <key>ProgramArguments</key>
  <array>
    <string>/Users/__username__/bin/homebrewupdate.sh</string> <!-- UPDATE file name -->
  </array>

  <key>Nice</key>
  <integer>1</integer>

  <key>StartInterval</key>
  <integer>21600</integer> <!-- UPDATE interval in seconds -->

  <key>RunAtLoad</key>
  <true/>

  <key>StandardErrorPath</key>
  <string>/tmp/__username__.homebrewupdate.err</string> <!-- UPDATE name/path -->

  <key>StandardOutPath</key>
  <string>/tmp/__username__.homebrewupdate.out</string><!-- UPDATE name/path -->
</dict>
</plist>
